# supervisord.conf.d/postgresql.conf
[program:postgresql]
command=/usr/lib/postgresql/11/bin/postgres -D /home/www-data/postgresql
stdout_logfile=/home/www-data/log/postgresql.log
redirect_stderr=true

