# supervisord.conf.d/tika.conf
[program:tika]
command=/usr/bin/java -jar /home/www-data/tika/tika-server.jar --host 127.0.0.1 --port 12345
stdout_logfile=/home/www-data/log/tika.log
redirect_stderr=true

